# ws2812b.service
[Unit]
Description=ws2812b.service
After=dnsmasq.service
[Service]
Type=forking
ExecStart=/home/pi/ws2812b.sh
[Install]
WantedBy=multi-user.target

